# com.sn2md.watch.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.sn2md.watch</string>

    <key>ProgramArguments</key>
    <array>
      <string>{{PROJECT_DIR}}/sn2md-watch.sh</string>
      <string>{{PROJECT_DIR}}/jobs.yaml</string>
    </array>

    <key>WorkingDirectory</key>
    <string>{{PROJECT_DIR}}</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>{{HOME}}/Library/Logs/sn2md-watch.log</string>
    <key>StandardErrorPath</key>
    <string>{{HOME}}/Library/Logs/sn2md-watch.err</string>
  </dict>
</plist>
